# service/thermia-publish.service
[Unit]
Description=Thermia Heatpump Avahi Publisher
After=network-online.target avahi-daemon.service
Wants=network-online.target
Requires=avahi-daemon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/thermia-publish.sh
Restart=always
RestartSec=10

NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
